# rootfs/etc/services.d/apache/run
#!/usr/bin/with-contenv sh

apache2-foreground
